# debian/octavia-driver-agent.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          octavia-driver-agent
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Should-Start:      postgresql mysql keystone rabbitmq-server ntp
# Should-Stop:       postgresql mysql keystone rabbitmq-server ntp
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Octavia LBaaS Driver agent
# Description:       Octavia Driver agent
### END INIT INFO

# Author: Frode Nordahl <frode.nordahl@canonical.com>
DESC="OpenStack Octavia Driver Agent (octavia-driver-agent)"
PROJECT_NAME=octavia
NAME=${PROJECT_NAME}-driver-agent
